# files/default/consul.path
[Unit]
Description=cloud conductor config file exists

[Path]
PathExists=/opt/cloudconductor/config
PathModified=/opt/cloudconductor/config
